# supervisor.conf
[supervisord]
nodaemon=true

[program:laravel]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/laravel.err.log
stdout_logfile=/var/log/laravel.out.log